# systemd/spopd.service
[Unit]
Description=Spop: a Spotify client daemon

[Service]
ExecStart=/usr/bin/spopd -f

[Install]
WantedBy=multi-user.target
